# deploy/nurtac-edge.service
[Unit]
Description=NurtacCoreEngineClaude Layer-14 Edge Matrix Engine
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/NurtacCoreEngineClaude
ExecStartPre=/bin/sleep 60
ExecStart=/root/NurtacCoreEngineClaude/.venv/bin/python3 edge_matrix_engine.py --mode live
Restart=always
RestartSec=5
MemoryMax=400M
MemorySwapMax=200M

[Install]
WantedBy=multi-user.target
