# deploy/nurtac-outcome.service
[Unit]
Description=NurtacCoreEngineClaude Layer-11 Historical Outcome Engine
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/NurtacCoreEngineClaude
ExecStartPre=/bin/sleep 45
ExecStart=/root/NurtacCoreEngineClaude/.venv/bin/python3 historical_outcome_engine.py --mode live
Restart=always
RestartSec=5
MemoryMax=400M
MemorySwapMax=200M

[Install]
WantedBy=multi-user.target
