# deploy/nurtac-observer.service
[Unit]
Description=NurtacCoreEngineClaude Layer-10 Observer Engine
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/NurtacCoreEngineClaude
ExecStartPre=/bin/sleep 40
ExecStart=/root/NurtacCoreEngineClaude/.venv/bin/python3 observer_engine.py --mode live
Restart=always
RestartSec=5
MemoryMax=400M
MemorySwapMax=200M

[Install]
WantedBy=multi-user.target
