# deploy/nurtac-layer2.service
[Unit]
Description=NurtacCoreEngineClaude Layer-2 (Aligned Candle Engine)
After=network-online.target nurtac-layer0.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/NurtacCoreEngineClaude
ExecStart=/root/NurtacCoreEngineClaude/.venv/bin/python3 aligned_candle_engine.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
