# deploy/nurtac-gate.service
[Unit]
Description=NurtacCoreEngineClaude Layer-5 Decision Gate
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/NurtacCoreEngineClaude
ExecStartPre=/bin/sleep 10
ExecStart=/root/NurtacCoreEngineClaude/.venv/bin/python3 decision_gate.py --mode live
Restart=always
RestartSec=5
MemoryMax=400M
MemorySwapMax=200M

[Install]
WantedBy=multi-user.target
