# deploy/nurtac-watchdog.service
[Unit]
Description=NurtacCoreEngineClaude Watchdog
After=network-online.target nurtac-layer0.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/NurtacCoreEngineClaude
ExecStart=/root/NurtacCoreEngineClaude/.venv/bin/python3 watchdog.py
Restart=always
RestartSec=10
# Optional: set Telegram credentials as environment variables
# Environment="TELEGRAM_BOT_TOKEN=<your-token>"
# Environment="TELEGRAM_CHAT_ID=<your-chat-id>"

[Install]
WantedBy=multi-user.target
