# deploy/nurtac-smartmoney.service
[Unit]
Description=NurtacCoreEngineClaude Layer-6 Smart Money Engine
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/NurtacCoreEngineClaude
ExecStart=/root/NurtacCoreEngineClaude/.venv/bin/python3 smart_money_engine.py --mode live
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
