# deploy/nurtac-volprofile.service
[Unit]
Description=NurtacCoreEngineClaude Layer-8 Volume Profile Engine
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/NurtacCoreEngineClaude
ExecStartPre=/bin/sleep 30
ExecStart=/root/NurtacCoreEngineClaude/.venv/bin/python3 volume_profile_engine.py --mode live
Restart=always
RestartSec=5
MemoryMax=400M
MemorySwapMax=200M

[Install]
WantedBy=multi-user.target
